# petalinux/fpgav3init/files/fpgav3init.service
[Unit]
Description=fpgav3init
 
[Service]
ExecStart=/usr/bin/fpgav3init
StandardOutput=journal+console
 
[Install]
WantedBy=multi-user.target
